# deploy/hive-autodeploy.service
[Unit]
Description=Hive pull-based auto-deploy (watch GHCR, redeploy on new image)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/root/projects/hive/scripts/auto-deploy.sh
